# systemd/astrolog-rs.service
[Unit]
Description=Astrolog-rs Server
After=network.target

[Service]
Type=simple
User=%i
Environment=PORT=4008
Environment=WORKERS=4
Environment=MAX_CONCURRENT=500
Environment=RUST_LOG=info
Environment=SE_EPHE_PATH=/home/%i/.swisseph/ephe
Environment=SE_LIB_PATH=/home/%i/.swisseph/lib
Environment=RUST_BACKTRACE=1

# System limits
LimitNOFILE=1048576
LimitNPROC=65535
LimitAS=infinity
LimitMEMLOCK=infinity
LimitSTACK=infinity
LimitCORE=infinity
LimitRSS=infinity
LimitLOCKS=infinity
LimitSIGPENDING=infinity
LimitMSGQUEUE=infinity
LimitNICE=infinity
LimitRTPRIO=infinity
LimitRTTIME=infinity

# TCP settings
Environment=TCP_NODELAY=1
Environment=TCP_QUICKACK=1
Environment=TCP_KEEPIDLE=60
Environment=TCP_KEEPINTVL=10
Environment=TCP_KEEPCNT=6

# Resource limits
MemoryLimit=4G
CPUQuota=90%

# Restart policy
Restart=always
RestartSec=5

# Working directory
WorkingDirectory=/home/%i/astrolog-rs

# Executable
ExecStart=/home/%i/astrolog-rs/target/release/astrolog-rs

[Install]
WantedBy=multi-user.target 